# archive/2026-06-06-wip/AdminPOOL/crontabs/README
In diesem Verzeichnis befinden sich die Skripte die
via crontab taeglich bzw. woechentlich ausgefuehrt
werden. Es gibt fuenf verschiedene Skripte:

1. cron.daily
	Macht ein CB_Make und CB_Make export
	(damit werden neu eingecheckte Dateien
	im Pool uebersetzt)

2. cron.weekly.clean
	Macht ein CB_Make clean
	(damit der ganze Pool nochmal uebersetzt wird,
	sollte kurz vor cron.daily gemacht werden)

3. cron.weekly.export
	Ruft CB_ExportProduct auf (kopieren der EXPORT-
	Verzeichnisse nach CB_Product, sollte kurz nach einem
	cron.daily gemacht werden)

4. a) cron.weekly.test
	Ruft CB_AutoTest aus CB_Admin/utils/CB_TestClient
	auf (startet ein paar mal den CBserver und fuehrt ein
	paar Skripte aus, wird kurz nach cron.weekly.export gemacht)

   b) cron.weekly.test.mail
    Checked, ob CB_AutoTest nicht-leere error.log-Files produziert
    und schickt diese an den ADMIN. Sollte gemacht werden, wenn
    cron.weekly.test auf allen Plattformen gelaufen ist.

5. cron.admin
    Ruft CB_Make clean, CB_Make, CB_Make export und CB_ExportAdmin im AdminPOOL auf,
    d.h. aktualisiert alle Programme und Dateien im AdminPOOL
    bzw. CB_Admin

6. daily.bat
    Ein Batch-File fuer Windows, das derzeit t輍lich um 2:25 Uhr
    auf Sargent laeuft (Scheduler ruft daily.bat auf), analog zu
    cron.daily.

7. autotest.bat
    Ein Batch-File, das den automatischen Test auf Windows durchf𢞴rt.


Diese Skripte sollten in $CB_ADMIN/crontabs zu finden sein.
Im Verzeichnis /home/cbase/crontabs gibt es Links zu diesen
Skripten. Dort werden auch die ganzen Log-Files und andere
Temp-Dateien gespeichert.

Die aktuelle crontab fuer cbase auf dix sieht so aus:
(weekly-scripts werden also mehrmals pro Woche ausgefuehrt)

10 0 1,16 * * /usr/ucb/mail -s "cron.error.log" cbase@i5.informatik.rwth-aachen.de < /home/CBase/crontabs/log/cron.error.log && /bin/rm /home/CBase/crontabs/log/cron.error.log && /bin/touch /home/CBase/crontabs/log/cron.error.log
15 0 * * 0,1,2,3,4,5,6 /usr/bin/date >> /home/CBase/crontabs/log/cron.error.log
30 0 * * 3 /home/CBase/crontabs/cron.weekly.clean 2>> /home/CBase/crontabs/log/cron.error.log
0 1 * * 0,1,2,3,4,5,6 /home/CBase/crontabs/cron.daily 2>> /home/CBase/crontabs/log/cron.error.log
0 3 * * 0,2,4,5 /home/CBase/crontabs/cron.weekly.export 2>> /home/CBase/crontabs/log/cron.error.log
0 4 * * 0,2,4,5 /home/CBase/crontabs/cron.weekly.test  2>> /home/CBase/crontabs/log/cron.error.log
44 4 * * 0,2,4,5 /home/CBase/crontabs/cron.weekly.test.mail  2>> /home/CBase/crontabs/log/cron.error.log
20 0 * * 2,5 /home/CBase/crontabs/cron.admin  2>> /home/CBase/crontabs/log/cron.error.log


Zeile 1: Zweimal im Monat (1. und 16.) wird das cron.error.log geschickt, und wieder auf Null gesetzt
Zeile 2: Taeglich wird im cron.error.log das Datum eingetragen, um Fehler bestimmten Tagen zu ordnen zu koennen
Zeile 3: Clean wird einmal pro Woche (Mittwochs) gemacht, vor allen anderen Jobs
Zeile 4: Daily wird wirklich taeglich durchgefuehrt, hier auf dix um 1 Uhr, auf den anderen Rechnern etwas spaeter
Zeile 5: Export wird Sonntags, Dienstags, Donnerstags und Freitags gemacht
Zeile 6: Der Test wird an den gleichen Tagen wie Export gemacht (Test benutzt die exportierten Daten)
Zeile 7: Nach dem Test muessen evtl. Fehler gemailt werden
Zeile 8: Der CB_Admin-Zweig wird zweimal pro Woche (dienstags und freitags) aktualisiert

Auf den anderen Plattformen wird nur Clean (Zeile 3), Daily (Zeile 4) und Test (Zeile 6) gemacht,
der Rest wird zentral nur auf sun4 (dix) gemacht.


