# AdminPOOL/crontabs/cron.weekly.test
#!/bin/csh
# This script runs several tests using CBshell
# Runs on all platforms (except Windows)

setenv PATH ${PATH}:/bin:/usr/bin

if ( "`uname`" == "Linux" ) then
        setenv PATH ${PATH}:/home/cbase/bin:/opt/sun-jdk-1.4.2.10/bin:/opt/blackdown-jdk-1.4.2.03/bin
else
        source /etc/dot/.environ
        source /etc/dot/.gui
        source /etc/dot/environments/012_cbase.env
endif

setenv CB_HOME /home/cbase/CB_NewStruct/CB_Product
setenv HOST `hostname`
setenv CB_POOL /home/cbase/CB_NewStruct/ProductPOOL
setenv CB_VARIANT `$CB_HOME/bin/CBvariant`
setenv LOG /home/cbase/crontabs/log/cron.weekly.test.$CB_VARIANT.log
setenv OLDLOG /home/cbase/crontabs/old
setenv SUBJECT "cron.weekly.test $HOST"
setenv ADMIN cbase@i5.informatik.rwth-aachen.de

if ( "$CB_VARIANT" == "sun4" ) then
   setenv MAILPROG /usr/ucb/mail
endif

if ( "$CB_VARIANT" == "i86pc" ) then
   setenv MAILPROG /usr/ucb/mail
   source /etc/dot/environments/009_jdk141_02.env
endif

if ( "$CB_VARIANT" == "linux" ) then
   setenv MAILPROG /usr/bin/mail
endif

/bin/mv -f $LOG $OLDLOG
touch $LOG

/bin/mv -f /home/cbase/crontabs/servertest/log/*.$CB_VARIANT /home/cbase/crontabs/servertest/old

cd /home/cbase/crontabs/servertest/$CB_VARIANT

echo "cron.weekly.test $HOST" `date` >> $LOG
echo . >> $LOG
echo "===================================* CB_AutoTest" >> $LOG
echo . >> $LOG

# kill existing CBserver processes of the current user (Aachen: cbase)
# could be left over by crashed executions of the 'weekly'/daily test scripts
pkill -u $USER CBserver

/home/cbase/CB_NewStruct/CB_Admin/utils/CB_TestClient/CB_AutoTest >>& $LOG

if ( $status != 0 ) then
	$MAILPROG -s "ERROR: $SUBJECT" $ADMIN < $LOG
else
	$MAILPROG -s "OK: $SUBJECT" $ADMIN < /dev/null
endif


