# AdminPOOL/crontabs/cron.weekly.test.mail
#!/bin/csh

# Script, which checks the error-files of the test procedure
# and sends mail to the admin, if they contain an error
# It should be run after all test procedures have been finished
# on all platforms.
# Runs only on Solaris SPARC

source /etc/dot/.environ
source /etc/dot/.gui

setenv ADMIN cbase@i5.informatik.rwth-aachen.de

cd /home/cbase/crontabs/servertest

foreach f  ( log/error.* )
	test -s $f
	if ( "$status" == "0" )  /usr/ucb/mail -s "$f" $ADMIN < $f
end

